# CONFIG/adad.service
[Unit]
Description=ada daemon

[Service]
ExecStart=/home/$USER/adasrc/bin
User=ada
Group=ada

[Install]
WantedBy=default.target
